# docker/quadlet/paperclip.container
[Unit]
Description=Paperclip AI Agent Orchestrator
Requires=offwatch-db.service
After=offwatch-db.service

[Container]
Image=offwatch-local
ContainerName=offwatch
Pod=offwatch.pod
Volume=%h/.local/share/offwatch:/offwatch:Z
Environment=HOST=0.0.0.0
Environment=PAPERCLIP_HOME=/offwatch
Environment=PAPERCLIP_DEPLOYMENT_MODE=authenticated
Environment=PAPERCLIP_DEPLOYMENT_EXPOSURE=private
Environment=PAPERCLIP_PUBLIC_URL=http://localhost:3100
EnvironmentFile=%h/.config/containers/systemd/offwatch.env

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
